# etc/services.d/postgres/finish
#!/bin/sh

S6_SERVICES=/var/run/s6/services
if s6-svok "$S6_SERVICES/postgres"; then 
 s6-svscanctl -t "$S6_SERVICES"
fi
