# etc/services.d/postgres/run
#!/usr/bin/with-contenv sh

source /etc/postgresql/common.sh

PGDATA=$DB_PATH
export PGDATA

exec s6-setuidgid postgres postgres -i 2>&1
